# crontab
PATH = '/home/ztatlock/cloth/bin:/usr/bin:/bin'

15 * * * * watchdog.sh
45 0 * * * analysis.sh
