# crontab
PATH   = '/home/ztatlock/cloth/bin:/home/ztatlock/bin:/usr/bin:/bin'
MAILTO = 'ztatlock@gmail.com'

19 * * * * /home/ztatlock/cloth/run.sh

