# crontab
PATH = '/usr/bin:/bin'

0 * * * * /home/ztatlock/cloth/bin/watchdog.sh
0 0 * * * /home/ztatlock/cloth/bin/daily.sh
