# crontab
PATH   = '/home/ztatlock/cloth:/home/ztatlock/bin:/usr/bin:/bin'
MAILTO = 'ztatlock@gmail.com'

19 * * * * /home/ztatlock/cloth/crawl.py

