# crontab
PATH = '/home/ztatlock/cloth/bin:/usr/bin:/bin'

0 * * * * watchdog.sh
0 0 * * * analysis.sh
