# crontab
PATH = '/home/ztatlock/clip/bin:/usr/bin:/bin'

15 *   * * * watchdog.sh
45 */4 * * * analysis.sh
